# scripts/com.exomem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd template rendered and installed by scripts/install-service.sh. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.exomem</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PYTHON__</string>
        <string>-m</string>
        <string>exomem</string>
        <string>--transport</string>
        <string>streamable-http</string>
        <string>--host</string>
        <string>__BIND_HOST__</string>
        <string>--port</string>
        <string>__PORT__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__WORKING_DIRECTORY__</string>

    __ENVIRONMENT_VARIABLES__

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/service.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/service.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
