# scripts/exomem.service
# systemd --user template for exomem. scripts/install-service.sh renders and
# installs this unit; do not install the template directly.

[Unit]
Description=exomem: Obsidian Knowledge Base MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=__SERVICE_ENV_FILE__
WorkingDirectory=__WORKING_DIRECTORY__
ExecStart="__VENV_PYTHON__" -m exomem --transport streamable-http --host __BIND_HOST__ --port __PORT__
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
